# recipes-connectivity/pi-ble-status/files/pi-ble-status.service
[Unit]
Description=Pi Status BLE GATT Server
After=bluetooth.target network.target
Requires=bluetooth.target

[Service]
Type=simple
ExecStart=/usr/bin/pi-ble-status
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
