# recipes-core/rauc-mark-good/files/rauc-mark-good.service
[Unit]
Description=RAUC Mark Current Slot as Good
Documentation=https://rauc.readthedocs.io/en/latest/integration.html
# Only run after the full multi-user target is up — this proves the system
# started successfully before we commit to this slot in U-Boot's BOOT_ORDER.
After=multi-user.target
Requires=multi-user.target
ConditionPathExists=/usr/bin/rauc

[Service]
Type=oneshot
ExecStart=/usr/bin/rauc status mark-good self
RemainAfterExit=yes
# Suppress error if already marked good (e.g., after a fresh flash of both slots)
SuccessExitStatus=0 1

[Install]
WantedBy=multi-user.target
